# infra/systemd/backend.service
[Unit]
Description=ClickDialog Backend (Gunicorn)
After=network.target

[Service]
User=www-data
WorkingDirectory=/opt/clickdialog/backend
Environment="FLASK_ENV=production"
EnvironmentFile=/opt/clickdialog/backend/.env
ExecStart=/opt/clickdialog/backend/.venv/bin/gunicorn -w 2 -b 127.0.0.1:5005 "app:create_app()"
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
